# container-fs/etc/cron.d/crontab
0 0 * * * /nightly.cron.sh >/dev/null 2>&1
*/15 * * * * /fifteen-minutes.cron.sh >/dev/null 2>&1
# An empty line is required at the end of this file for a valid cron file.